# displaylink-gentoo/init.d/displaylink
#!/sbin/openrc-run

name="DisplayLinkManager"
command="/opt/displaylink/DisplayLinkManager"
pidfile="/run/displaylink.pid"
command_background=true
description="DisplayLink USB graphics driver"
